# stack/config/jenkins/system/profile.d/nodejs.sh
export NODE_HOME=/usr/share/nodejs/node-v11.10.0-linux-x64

export PATH=${PATH}:${NODE_HOME}/bin
